# deploy.conf
[program:goms]
command=/opt/goms/goms
autostart=true
autorestart=true
startsecs=5
startretries=3
